# config/systemd/remotehive-backend.service
[Unit]
Description=RemoteHive Backend API Service
After=network.target mongodb.service redis.service
Wants=mongodb.service redis.service
Requires=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/RemoteHive
Environment=PATH=/home/ubuntu/.local/bin:/usr/bin:/bin
EnvironmentFile=/home/ubuntu/RemoteHive/.env
ExecStart=/usr/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 2
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/ubuntu/RemoteHive/logs /tmp /home/ubuntu/RemoteHive/uploads
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Health check
ExecStartPost=/bin/sleep 10
ExecStartPost=/bin/bash -c 'curl -f http://localhost:8000/health || exit 1'

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=remotehive-backend

[Install]
WantedBy=multi-user.target